# etc/systemd/user/update-user-flatpaks.service
[Unit]
Description=Update user Flatpaks

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update --assumeyes --noninteractive

[Install]
WantedBy=default.target
